# samples/sample.cron
* * * * * echo "cron out"
* 1,2,3,4,5 * * * echo "cron out"
